# virtual-office.service
[Unit]
Description=Virtual Office — Backend + Frontend
After=network.target
Wants=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/projects/virtual-office-poc/server
EnvironmentFile=%h/projects/virtual-office-poc/server/.env
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=STATIC_DIR=%h/projects/virtual-office-poc/dist
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
